# run.cron
@reboot cd /home/mo/code/wifispeed && /usr/bin/python3 /home/mo/code/wifispeed/wifispeed.py &
